# systemd-file.service
[Unit]
Description=Your Application Description
After=network.target

[Service]
User=csye6225
Group=csye6225
WorkingDirectory=/path/to/your/application/
ExecStart=/path/to/your/application/start-command
Restart=always

[Install]
WantedBy=multi-user.target
